# config/move_footage.service
[Unit]
Description=Moves video footage generated by the motion detection software from this machine to a nas server. 

[Service]
Type=oneshot
ExecStart=/home/weather/weather/util/move_footage.py
User=weather
