# config/move_footage.timer
[Unit]
Description=Run the move_footage.service every 2 minutes

[Timer]
OnCalendar=*:0/2
